# config/supervisor/conf.d/event_server.conf
[program:event_server]
command=/opt/wxovercast/bin/event_server
process_name=%(program_name)s
autostart=true
autorestart=true
startsecs=5
startretries=3
environment = PYTHONPATH="/opt/wxovercast"
